# src/broker/tme-spyd.cron
*/2 * * * * root /usr/share/mist/bin/watchdog-spyd >> /var/run/tme/tme-spyd.watch
