# src/zookeeper/etc/tme-zookeeper.cron
* * * * * root /usr/share/tme-zookeeper/watchdog-zookeeperd >> /var/run/tme/zookeeper.watch
