# src/mist/build/tme-bridge
#!/bin/sh
# chkconfig: 345 25 75
# description: TME bridge daemon 

prog="tme-bridge"
pidfile="/var/run/tme/$prog.pid"
errfile="/var/run/tme/$prog.err"
lockfile="/var/lock/subsys/$prog"
ret_code=0
classname='com.trendmicro.mist.cmd.TmeBridge'

javacmd="/usr/java/latest/bin/java"
if [ ! -x "$javacmd" ]; then
    echo "Sun Java Runtime Environment: $javacmd not found"
    exit 1
fi

check_process() {
    pattern=$1
    try_count=$2
    exist=$3        # 1: check exist
                    # 0: check not exist
    counter=0
    ready=0
    while [ $counter -lt $try_count ] ; do
        pid=`ps -C java -o pid,args | grep $pattern | sed -e 's/^ *//' | cut -d ' ' -f 1`
        if [ $exist -eq 1 ] && [ ! -z "$pid" ]; then
            ready=1
            break
        elif [ $exist -eq 0 ] && [ -z "$pid" ]; then
            ready=1
            break
        else
            usleep 50000
            let counter=counter+1
        fi
    done

    if [ $ready -eq 1 ]; then
        return 0
    else
        return 1
    fi
}

check_pid() {
    my_pid=$1
    if [ -e $my_pid ]; then
        ps -p `cat $my_pid` > /dev/null
        return $?
    else
        return 1
    fi
}

start_daemon() {
    echo -n "Starting $prog ... "
    check_pid $pidfile
    running=$?
    if [ $running -eq 0 ]; then
        echo "already started"
    else
        su TME -c "nohup $javacmd -cp \"/usr/share/mist/lib/*:/usr/share/spn-infra-common/lib/*\" $classname >> $errfile 2>&1 &"
        check_process $classname 60 1
        ready=$?
        if [ $ready -eq 0 ]; then
            echo "success"
            touch $lockfile
        else
            echo "failed"
            cat $errfile
	        ret_code=1
        fi
    fi
}

stop_daemon() {
    echo -n "Stoping $prog ... "
    check_pid $pidfile
    running=$?
    if [ $running -eq 0 ]; then
        kill `cat $pidfile`
        check_process $classname 60 0
        ready=$?
        if [ $ready -eq 0 ]; then
            echo "success"
        else
            echo "failed, force kill"
            kill -9 `cat $pidfile`
            rm -f $pidfile
        fi
        rm -f $lockfile
    else
        echo "not started"
    fi
}

case "$1" in
start)
    start_daemon
    ;;
stop)
    stop_daemon
    ;;
restart)
    stop_daemon
    start_daemon
    ;;
status)
    check_process $classname 5 1
    ret_code=$?
    if [ $ret_code -eq 0 ]; then
	    echo "deamon is running"
    else
	    echo "deamon is stopped"
    fi
    ;;
*)
    echo "Usage: $0 {start|stop|restart|status}"
    ;;
esac

exit $ret_code
