# src/mist/etc/tme-bridge.cron
*/2 * * * * root /usr/share/mist/bin/watchdog-bridge >> /var/run/tme/tme-bridge.watch
