# src/mist/etc/tme-mist.cron
*/5 * * * * root /usr/share/mist/bin/watchdog-mistd >> /var/run/tme/mistd.watch
